# crontab
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=$USER
HOME=/home/$USER

0  0  *  *  * /var/www/django/jaxid_generator/check_and_backup_jaxid_db.sh
